# crontab
# local user crontab file

0 0 * * * find ${HOME}/workspace/ -name development.log -exec cp /dev/null {} \;
