# crontab
# local user crontab file

0 0 * * * find /Users/robbycolvin/workspace/ -name development.log -exec cp /dev/null {} \;
